# systemd.service
[Unit]
Description=Bussiaeg.ee
After=mysql.service

[Service]
ExecStart=/bin/bash start.sh

[Install]
WantedBy=multi-user.target
